# socat.conf
[program:socat]
command=/opt/bin/socat.sh
autostart=true
autorestart=true
startretries=3
startsecs=15
stopwaitsecs=30

redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
